# infra/systemd/wisp-ops-web.service
[Unit]
Description=wisp-ops-center Web (Next.js)
After=network.target

[Service]
Type=simple
User=wispops
Group=wispops
WorkingDirectory=/opt/wisp-ops-center/apps/web
EnvironmentFile=/etc/wisp-ops/web.env
ExecStart=/usr/bin/node /opt/wisp-ops-center/apps/web/node_modules/next/dist/bin/next start -p 3000
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
