# infra/systemd/wisp-ops-worker.service
[Unit]
Description=wisp-ops-center Worker
After=network.target postgresql.service redis-server.service
Wants=postgresql.service redis-server.service

[Service]
Type=simple
User=wispops
Group=wispops
WorkingDirectory=/opt/wisp-ops-center
EnvironmentFile=/etc/wisp-ops/worker.env
ExecStart=/opt/wisp-ops-center/bin/wisp-ops-worker
Restart=on-failure
RestartSec=10
LimitNOFILE=65536
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/var/log/wisp-ops
ProtectHome=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
